# systemd/license-sync.timer
[Unit]
Description=Run LICENSE sync weekly

[Timer]
OnCalendar=Sun *-*-* 04:47:00
Persistent=true

[Install]
WantedBy=timers.target
