# systemd/readme-updater.timer
[Unit]
Description=Run README updater daily

[Timer]
OnCalendar=*-*-* 04:17:00
Persistent=true

[Install]
WantedBy=timers.target
